# com.jason.claude-plan-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jason.claude-plan-usage</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Jason/.volta/bin/node</string>
        <string>/Users/Jason/Development/claude-statusline/dist/claude-plan-usage.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Jason/Development/claude-statusline</string>

    <key>StartInterval</key>
    <!-- Run every 10 minutes (600 seconds) -->
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/Users/Jason/Library/Logs/claude-plan-usage.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Jason/Library/Logs/claude-plan-usage.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Jason/.volta/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
