# com.username.claude-plan-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
    Claude Plan Usage Scraper - launchd Job

    Installation:
    1. Copy this file to ~/Library/LaunchAgents/
    2. Edit paths below to match your setup
    3. Run: launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.username.claude-plan-usage.plist

    To unload: launchctl bootout gui/$(id -u)/com.username.claude-plan-usage
    To run now: launchctl kickstart -k gui/$(id -u)/com.username.claude-plan-usage
    -->

    <key>Label</key>
    <string>com.username.claude-plan-usage</string>

    <key>ProgramArguments</key>
    <array>
        <!-- EDIT: Path to your node binary (run 'which node' to find it) -->
        <string>/opt/homebrew/bin/node</string>
        <!-- EDIT: Full path to the compiled script -->
        <string>/Users/username/Development/claude-statusline/dist/claude-plan-usage.js</string>
        <string>--profile</string>
        <string>Profile 7</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- EDIT: Path to the project directory -->
    <string>/Users/username/Development/claude-statusline</string>

    <key>StartInterval</key>
    <!-- Run every 10 minutes (600 seconds) -->
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/Users/username/Library/Logs/claude-plan-usage.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/username/Library/Logs/claude-plan-usage.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Ensure node can find modules -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>

        <!-- Chrome profile name (Profile 7 = "Person 1") -->
        <key>CHROME_PROFILE_NAME</key>
        <string>Profile 7</string>
    </dict>

    <!-- Run at load (useful for testing) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Keep trying if it exits non-zero -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Limit retries to avoid hammering on persistent failures -->
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <!-- Nice level: run at lower priority -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
